# fboss-image/image_builder/templates/centos-09.0/root_files/usr/lib/systemd/system/fsdb.service
# /usr/lib/systemd/system/fsdb.service
[Unit]
Description=FBOSS FSDB Service
After=rc-local.service
StartLimitIntervalSec=600
StartLimitBurst=5

[Service]
Type=simple
LimitNOFILE=65535
LimitCORE=32G

Environment="PATH=/opt/fboss/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin"
Environment="LD_LIBRARY_PATH=/opt/fboss/lib"
WorkingDirectory=/opt/fboss
ExecStart=/bin/bash -c 'source /opt/fboss/bin/setup_fboss_env && exec /opt/fboss/bin/fsdb'

Restart=always
RestartSec=30
TimeoutStartSec=180s

# Logging to journald
StandardOutput=journal
StandardError=journal
SyslogIdentifier=fsdb

[Install]
WantedBy=multi-user.target
