# fboss-image/image_builder/templates/centos-09.0/root_files/usr/lib/systemd/system/local_rpm_repo.service
[Unit]
Description=Enable local_rpm_repo

[Service]
Type=oneshot
ExecStart=/usr/bin/createrepo /usr/local/share/local_rpm_repo
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
