# fboss-image/image_builder/templates/centos-09.0/root_files/etc/profile.d/fboss.sh
#!/bin/bash
# shellcheck source=/dev/null
cd /opt/fboss && source /opt/fboss/bin/setup_fboss_env
cd - > /dev/null || return
